# deploy/systemd/adreem-api.service
[Unit]
Description=ADREEM Ledger API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/argaz/apps/adreem
EnvironmentFile=/home/argaz/apps/adreem/adreem.env
ExecStart=/usr/bin/npm run api:adreem
Restart=always
RestartSec=5
StandardOutput=append:/home/argaz/logs/adreem-api.log
StandardError=append:/home/argaz/logs/adreem-api-error.log

[Install]
WantedBy=default.target
